# build/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid
user=root
# 优雅退出配置
stopwaitsecs=30
# 最小可用文件描述符
minfds=1024

; XML-RPC server for services plugin (Story 6.1)
; Note: For development, using simple password. In production, use ENV variable.
[inet_http_server]
port=127.0.0.1:9001
username=admin
password=%(ENV_SUPERVISOR_PASSWORD)s

[supervisorctl]
serverurl=http://127.0.0.1:9001
username=admin
password=%(ENV_SUPERVISOR_PASSWORD)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ======================
; 系统基础服务
; ======================

[program:redis]
command=/usr/bin/redis-server --appendonly yes --dir /appos/data/redis
priority=10
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
user=root
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

[program:pocketbase]
command=/usr/local/bin/pocketbase serve --dir /appos/data/pocketbase --http 127.0.0.1:8090
priority=15
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
user=root
stdout_logfile=/var/log/supervisor/pocketbase.log
stderr_logfile=/var/log/supervisor/pocketbase-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

; ======================
; 业务服务
; ======================

[program:backend]
command=/usr/local/bin/appos-api
directory=/appos/data
environment=HOME="/appos/data",REDIS_ADDR="127.0.0.1:6379",POCKETBASE_URL="http://127.0.0.1:8090"
priority=20
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=15
user=root
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=QUIT
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
